# infra/podman/quadlet/reseau-web.container
[Unit]
Description=reseau-web (dev)

[Container]
ContainerName=reseau-web
Image=localhost/reseau-web:dev
Network=host
Environment=NEXT_PUBLIC_API_BASE=http://localhost:8000
Volume=%h/projects/reseau-social/web:/app/web:Z
Exec=bash -lc "cd /app/web && npm run dev"

[Install]
WantedBy=default.target
